# sample/crontab
5  */1  *   *   *   root  { cd /opt/nfs_sync && nfs_check.sh && DIR=mydir nfs_sync.sh; } > /var/log/nfs_sync/nfs_sync_mydir.log 2>&1

#0   23   *  *   sat   root  { cd /opt/nfs_sync && nfs_check.sh && DIR=mydir && nfs_fullsync.sh; } > /var/log/nfs_sync/nfs_fullsync_mydir.log 2>&1
